# dockerify/uwsgi/mysana_celery.conf
start on filesystem or runlevel [2345]
stop on shutdown or runlevel [!2345]
respawn

sudo /etc/init.d/mysana_celeryd start